# sources/helper/netctlgui-helper.service
[Unit]
Description=netctlgui-helper daemon

[Service]
ExecStart=/usr/bin/netctlgui-helper

[Install]
WantedBy=multi-user.target
